# conf/rc.local
#!/bin/sh -e

service mysql start
service nginx start
python3 $root_dir/$app/manage.py migrate
echo yes | python3 $root_dir/$app/manage.py collectstatic
cd $root_dir
gunicorn $app.wsgi:application
